# archlinux/mkosi.extra/etc/systemd/system/getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --noreset --noclear --nohostname --login-options "-H" - ${TERM}
