# archlinux/mkosi.extra/etc/profile.d/90-shell.sh
if [ -n "${BASH_VERSION-}" ]; then
    PS1='[\u@\H \W]\$ '
fi
